# renew-cert.service
[Unit]
Description=Update thefiletree.com TLS certificate

[Service]
Type=oneshot
ExecStart=/home/dom/services/renew-cert
